# crontab/py_job.sh
#!/bin/bash

/Users/charilaostsarouchas/PT_local/cronjob/py_job.py
