# crontab/spark_pi.sh
#!/bin/bash
/Users/charilaostsarouchas/tsarouch_local/work_soft/spark-1.5.2-bin-hadoop2.6/bin/spark-submit /Users/charilaostsarouchas/PT_local/cronjob/spark_pi.py
#/Users/charilaostsarouchas/PT_local/cronjob/spark_pi.py
